# oacisjupyterlab/jupyter.conf
[program:jupyter]
command=bash -l -c 'pipenv run jupyter lab'
directory=/home/oacis
user=oacis
autostart=true
autorestart=true
environment=SHELL="/bin/bash"
